# workspace/PrintPriceOS_Workspace/deployment/ppos-preflight-worker.service
[Unit]
Description=PrintPrice OS Preflight Worker
After=network.target redis.target

[Service]
Type=simple
User=ppos-admin
WorkingDirectory=/opt/printprice-os/ppos-preflight-worker
Environment=NODE_ENV=production
Environment=HEALTH_PORT=8002
Environment=PPOS_UPLOADS_DIR=/opt/printprice-os/temp-staging
ExecStart=/usr/bin/node worker.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
