# workspace/PrintPriceOS_Workspace/deployment/ppos-preflight-service.service
[Unit]
Description=PrintPrice OS Preflight Service
After=network.target redis.target

[Service]
Type=simple
User=ppos-admin
WorkingDirectory=/opt/printprice-os/ppos-preflight-service
Environment=NODE_ENV=production
Environment=PPOS_SERVICE_PORT=8001
Environment=PPOS_UPLOADS_DIR=/opt/printprice-os/temp-staging
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
